# cron/crontab
* * * * * c:\Server\bin\PHP c:\Server\data\htdocs\ImportantProjects\cron\cron_app\index.php